# features/container-host/rootfs/usr/lib/systemd/system/podman-auto-update.timer.d/10-schedule.conf
[Timer]
# Check for updates daily at 4am
OnCalendar=
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=30min
